# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[fcgi-program:uvicorn]
socket=unix:///tmp/uvicorn.sock
socket_owner=www-data
directory=/sgs_caminho_critico
command=uvicorn --fd 0 run:app
numprocs=5
process_name=uvicorn-%(process_num)d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
capture_output = True
errorlog = '/tmp/uvicorn.error.log'
accesslog = '/tmp/log/uvicorn.log'
